# utils/nctl-metrics/supervisord.conf
[supervisord]
nodaemon = True

[supervisorctl]
serverurl = unix://%(here)s/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=%(here)s/supervisord.sock

[program:memory-stats-collector]
directory=%(here)s
command=%(here)s/mem_export.py
stdout_logfile=%(here)s/memory-stats-collector.log
redirect_stderr=True
environment=PYTHONUNBUFFERED=1

[program:prometheus-container]
directory=%(here)s
command=%(here)s/prometheus.sh
stdout_logfile=%(here)s/prometheus.log
redirect_stderr=True
